# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/data/supervisord.pid

[program:rlm]
command=rlm-mcp
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/rlm.log

[program:jefest]
command=python -m jefest.server
directory=/app
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/jefest.log
